# deployment/supervisor-rescuelink-worker.conf
[program:rescuelink-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/rescuelink/backend/artisan queue:work database --sleep=1 --tries=3 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/rescuelink/backend/storage/logs/worker.log
stopwaitsecs=30
